# picar.service
[Unit]
Description=Pycar Project
After=Network.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/picar.py

[Install]
WantedBy=multi-user.target
